# baseimg/storpool/storpool-18.02.752.1736934/mgmt/etc/cron.d/storpool_monitor
# /etc/cron.d/storpool_storpool_monitor: crontab entries for the storpool monitoring pusher

PATH=/bin:/usr/sbin:/usr/bin

* *    * * *   root [ -f /etc/storpool.conf ] && /usr/lib/storpool/monitor
