# baseimg/storpool/storpool-18.02.752.1736934/common/etc/cron.d/storpool_setirqaff
# /etc/cron.d/storpool_setirqaff: crontab entries for the StorPool Set IRQ Affinity tool

PATH=/bin:/sbin/:/usr/sbin:/usr/bin

*/5 *    * * *   root  /usr/lib/storpool/setirqaff do log >/dev/null 2>&1
