# baseimg/storpool/storpool-18.02.752.1736934/mgmt/usr/share/doc/storpool/examples/etc/cron.d/storpool_scrub_disks
#
#-
# Copyright (c) 2016  StorPool.
# All rights reserved.
#

# start at midnight on Saturdays
0 0 * * 6 root  [ -f /etc/storpool.conf ] && /usr/lib/storpool/scrub_disks -a --start
# example stop scrubbing at 8AM on Sundays
#0 8 * * 0 root [ -f /etc/storpool.conf ] && /usr/lib/storpool/scrub_disks -a --stop
