# baseimg/storpool/storpool-18.02.752.1736934/repsync/usr/lib/storpool/storpool_repsync.cron
# /etc/cron.d/storpool_repsync: crontab entries for the StorPool report sync tool

PATH=/bin:/usr/sbin:/usr/bin

* *    * * *   %%USERNAME%% [ -f /etc/storpool.conf ] && /usr/sbin/storpool_repsync
