# baseimg/storpool/storpool-18.02.752.1736934/common/usr/lib/systemd/system/storpool.slice
#
#-
# Copyright (c) 2013  StorPool.
# All rights reserved.
#

[Unit]
Description=Storpool Slice
DefaultDependencies=no
Before=slices.target
Wants=-.slice
After=-.slice

[Slice]
MemoryAccounting=true
