# mdbbackup_init/periodic/15min/test
#!/bin/sh
NOW=$(date '+%M%H%d')
BDIR=/backup/hourly/
touch "$BDIR"/test-horly-"$NOW"

